# bin_hexapod/DEBIAN/postinst
#!/bin/bash

update-rc.d hexapod defaults

echo "enjoy your hexapod ! ;)"