# bin_hexapod/DEBIAN/prerm
#!/bin/bash

update-rc.d -f hexapod remove